# .config/systemd/user/borgbackup.service
[Unit]
Description=BorgBackup
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c '~/bin/bb'

[Install]
WantedBy=default.target
