# .config/systemd/user/borgbackup.timer
[Unit]
Description=BorgBackup

[Timer]
OnCalendar=*-*-* 17:30:00

[Install]
WantedBy=timers.target
